# sheerdns.service
[Unit]
Description=sheerdns service
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/sheerdns

[Install]
WantedBy=multi-user.target
